# node-red.conf
start on started mountall
stop on shutdown

script
  exec /usr/local/bin/node /usr/local/bin/node-red > /var/log/node-red.log 2>&1
end script
